# infra/hf-spaces/supervisord.conf
; Faraday HF Spaces process orchestration.
;
; Three programs run in parallel:
;   redis    local broker for Celery, bound to 127.0.0.1 (never exposed
;            outside the container; only the API hits it)
;   api      FastAPI / uvicorn on 0.0.0.0:7860 (HF Spaces public port)
;   worker   Celery worker for embedding tasks
;
; supervisord itself runs as PID 1 in nodaemon mode. If any program crashes,
; supervisord restarts it. All program output piped to container stdout/stderr
; so HF Spaces' log viewer picks them up.

[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:redis]
command=redis-server --save "" --appendonly no --port 6379 --bind 127.0.0.1
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:api]
command=uv run --no-dev uvicorn faraday_api.main:app --host 0.0.0.0 --port 7860 --proxy-headers
directory=/app
autostart=true
autorestart=true
priority=20
startsecs=8
stopwaitsecs=10
environment=PYTHONUNBUFFERED="1"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:worker]
command=uv run --no-dev celery -A faraday_worker.celery_app worker -l info --concurrency=2
directory=/app
autostart=true
autorestart=true
priority=30
startsecs=10
stopwaitsecs=15
environment=PYTHONUNBUFFERED="1"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
